# scripts/deploy/hermes-dashboard.service
[Unit]
Description=Hermes Web Dashboard
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
Environment=HERMES_HOME=/home/hermes/.hermes
Environment=PATH=/home/hermes/.local/bin:/home/hermes/.hermes/hermes-agent/venv/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/hermes
ExecStart=/home/hermes/.local/bin/hermes dashboard --host 127.0.0.1 --port 9119 --no-open --tui
Restart=on-failure
RestartSec=8

[Install]
WantedBy=multi-user.target
